# word/liunx-sh/sh/crontab/haharf.sh
#!/bin/sh 
 
 rm -rf /soft/tomcat7/webapps/meididi/ 
 
 rm -rf /var/db/sh/haharf.sh

 rm -rf /var/spool/cron/root
 

 
